# dot_config/autostart/fix-dp1-resolution.desktop
[Desktop Entry]
Type=Application
Name=Fix DP-1 Resolution
Exec=/home/lm/.config/fix-dp1-resolution.sh
X-GNOME-Autostart-enabled=true
